# debian/pve-san-fenced.service
[Unit]
Description=Proxmox VE SAN Fencing Daemon
Documentation=man:pve-san-fenced(8)
After=network.target multipathd.service
Requires=multipathd.service

[Service]
Type=simple
EnvironmentFile=-/etc/default/pve-san-fenced
ExecStart=/usr/sbin/pve-san-fenced
Restart=always
RestartSec=5
TimeoutStartSec=120
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
RuntimeDirectory=pve-san-fenced
ReadWritePaths=/proc/sysrq-trigger

[Install]
WantedBy=multi-user.target
